# letsencrypt_acme/rootfs/etc/services.d/lets-encrypt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start sshd service if enabled
# ==============================================================================
CERT_DIR=/data/letsencrypt
WORK_DIR=/data/workdir
PROVIDER_ARGUMENTS=()
ACME_CUSTOM_SERVER_ARGUMENTS=()

EMAIL=$(bashio::config 'email')
DOMAINS=$(bashio::config 'domains')
KEYFILE=$(bashio::config 'keyfile')
CERTFILE=$(bashio::config 'certfile')
ACME_SERVER=$(bashio::config 'acme_server')
ACME_ROOT_CA=$(bashio::config 'acme_root_ca_cert')
IONOS_PREFIX=$(bashio::config 'ionos_prefix')
IONOS_SECRET=$(bashio::config 'ionos_secret')

export IONOS_PREFIX=$IONOS_PREFIX
export IONOS_SECRET=$IONOS_SECRET

/root/.acme.sh/acme.sh --register-account -m $EMAIL
/root/.acme.sh/acme.sh --issue --dns dns_ionos -d $DOMAINS --key-file /ssl/$KEYFILE --cert-file /ssl/$CERTFILE